# systemd_files/usb-daemon.service
[Unit]
Description=USB demon service
DefaultDependencies=no

[Service]
ExecStart=/bin/usb_demon

[Install]
WantedBy=usb-gadget.target
