# stage/etc/systemd/system/usb-daemon.service
[Unit]
Description=USB demon service
DefaultDependencies=no

[Service]
ExecStart=/bin/usb_demon
